# anouman/templates/init_script/website.conf
description     "Starting {{DOMAINNAME}}"

start on startup
exec {{GUNICORN_START}}
respawn
